# mfgflow.service
[Unit]
Description=MFGFLOW Production Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/mfgflow
EnvironmentFile=/home/pi/mfgflow/mfgflow.env
ExecStart=/usr/bin/python3 /home/pi/mfgflow/pi_server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
